# vacuum_supervisor/robotcontroller.conf
[inet_http_server]
port=5000

[program:controller]
#command=%(here)s/run_with_env.sh /home/pi/vacuumws/devel/setup.bash rosrun vacuum_controllers robot_controller.py --config cfg/%(host_node_name)s.cfg
command=/home/pi/Source/vacuum_supervisor/run_with_env.sh /home/pi/vacuumws/devel/setup.bash rosrun vacuum_controllers robot_controller.py --config cfg/%(host_node_name)s.cfg
user=pi
stopsignal=INT
stdout_logfile=/home/pi/logs/controller.log
redirect_stderr=true
#shell=true
autostart=false
autorestart=false
